# files/etc/systemd/system/hakko-pull.timer
[Unit]
Description=Run Hakko pull periodically

[Timer]
OnBootSec=20
OnUnitActiveSec=60

[Install]
WantedBy=timers.target

