# files/etc/systemd/system/hakko-push.service
[Unit]
Description=Push files from local in directory to Hakko FTP

[Service]
Type=oneshot
ExecStart=/usr/local/bin/hakko-push.sh

